# upstart.conf

start on started networking

respawn

setuid ubuntu
setgid ubuntu

chdir /home/ubuntu/colour-demo

script
    /usr/bin/python colour.py
end script
